# debian/j1708-drivers/etc/init/ecm.conf
start on start_j1708
stop on stop_j1708

exec /opt/synercon/bin/ecm

post-start script
	   initctl emit ecm-started
end script