# debian/j1708-drivers/etc/init/fwd-j1708.conf
start on (start_j1708 and (started dpa and started ecm))
stop on stop_j1708

exec /opt/synercon/bin/J1708PassthroughDriver

